# supervisor/movie.ini
[program:wechat-movie]
command=/data/virtualenv/wechat-movie/bin/python bin/index.py
numprocs=1
directory=/data/wechat-movie
autostart=true
autorestart=true
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
