# supervisor/movie.ini
[program:wechat_movie]
command=/data/wechat-movie/venv/bin/python server.py
process_name=%(program_name)s
numprocs=1
directory=/data/wechat-movie
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/wechat_movie/wechat_movie.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/wechat_movie/wechat_movie.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
enviroment=PYTHONPATH='/data/wechat-movie'

[program:wechat_movie.celery.worker]
command=/data/wechat-movie/venv/bin/celery -A movie.celery_tasks worker -l debug -n wechat_movie_celery_worker%(process_num)02d.%%h --heartbeat-interval=30
process_name=%(program_name)s
numprocs=1
directory=/data/wechat-movie
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/wechat_movie/wechat_movie_celery_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/wechat_movie/wechat_movie_celery_worker.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
enviroment=PYTHONPATH='/data/wechat-movie'

[program:wechat_movie.celery.beat]
command=/data/wechat-movie/venv/bin/celery -A movie.celery_tasks beat -l debug -s /data/log/wechat_movie/celerybeat.schedule.cache --pidfile=/data/log/wechat_movie/celerybeat.pid.cache
process_name=%(program_name)s
numprocs=1
directory=/data/wechat-movie
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/wechat_movie/wechat_movie_celery_beat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/wechat_movie/wechat_movie_celery_beat.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
enviroment=PYTHONPATH='/data/wechat-movie'
